# unlock-agent/systemd/sailfish-unlock-agent.service
[Unit]
Description=Sailfish Unlock Agent
DefaultDependencies=no
Conflicts=shutdown.target emergency.service display.service
Before=shutdown.target display.service
ConditionDirectoryNotEmpty=/run/systemd/ask-password

[Service]
Type=oneshot
ExecStart=/usr/libexec/sailfish-unlock-ui
