# unlock-agent/systemd/sailfish-unlock-agent.service
[Unit]
Description=Sailfish Unlock Agent
DefaultDependencies=no
Conflicts=shutdown.target emergency.service display.service
Before=shutdown.target display.service

[Service]
ExecStart=/usr/bin/unlock-agent
SuccessExitStatus=0 8
