# unlock-agent/systemd/systemd-cryptsetup@.service.d/50-actdead.conf
# Skip unlocking in actdead mode
[Unit]
After=initial-bootstate.service
ConditionPathExists=!/run/systemd/boot-status/ACT_DEAD
