# unlock-agent/systemd/systemd-cryptsetup@.service.d/50-factorytest.conf
# Skip unlocking in factory test mode
[Unit]
After=initial-bootstate.service
ConditionPathExists=!/run/systemd/boot-status/TEST
